# crontab
# run backup every night at 2 AM
* * * * * node /app/backup.js backup >> /proc/1/fd/1 2>&1
