# crontab
0 2 * * * node /app/backup.js backup >> /proc/1/fd/1 2>&1